# roles/service/templates/traxconsole.service.j2
# container-traxconsole.service
# autogenerated by Podman 1.9.3
# Wed Jul 29 14:25:46 CDT 2020

[Unit]
Description=Podman container-traxconsole.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStart=/usr/local/bin/startTraxConsole.sh
ExecStop=/usr/local/bin/stopTraxConsole.sh
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
